# roles/tor/files/supervisor-tor.conf
[program:tor]
command=/usr/sbin/tor -f /etc/tor/torrc
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
